# static/etc/services.d/murmur/run
#!/usr/bin/with-contenv bash

umask 0002

echo "Starting Murmur"
exec s6-setuidgid bytesized /app/murmur.x86 -ini /config/murmur.ini -fg
